# LinuxInstall.SettingFiles/elm327_websocket.conf
[program:elm327_websocket]
command=mono ELM327COM_WebSocket_Server.exe
user=wscomm
directory=/home/wscomm/WebsocketServer/ELM327
autorestart=true
stdout_logfile=/var/log/supervisor/elm327_websocket.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
redirect_stderr=true
